# com.officetracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.officetracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/rahulmishra/Desktop/Personal/wifi-tracking/venv/bin/python</string>
        <string>-m</string>
        <string>uvicorn</string>
        <string>app.main:app</string>
        <string>--host</string>
        <string>127.0.0.1</string>
        <string>--port</string>
        <string>8787</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rahulmishra/Desktop/Personal/wifi-tracking</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/rahulmishra/Desktop/Personal/wifi-tracking/logs/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rahulmishra/Desktop/Personal/wifi-tracking/logs/stderr.log</string>
</dict>
</plist>
